# root/etc/s6-overlay/s6-rc.d/svc-coturn/run
#!/usr/bin/with-contenv bash
set -e

COTURN_ARGS="${COTURN_ARGS:- -v}";

s6-setuidgid ${S6_USER:-alpine} exec \
    turnserver \
    -c /var/lib/coturn/turnserver.conf \
    ${COTURN_ARGS}
